# infra/chess-time-frontend.service
[Unit]
Description=Chess Time Next.js Frontend
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/chess-time-frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/chess-time-frontend/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
